# patches/vban-manager/vban@.service.in
[Unit]
Description=VBAN stream %i
After=default.target

[Service]
Type=simple
WorkingDirectory=@MANAGER_DIR@/script
ExecStart=/bin/bash @MANAGER_DIR@/script/vban.sh start-service %i
Restart=on-failure
RestartSec=3
StartLimitIntervalSec=120
StartLimitBurst=5

[Install]
WantedBy=default.target
